# deploy/mutsa-delivery.service
[Unit]
Description=Mutsa Delivery API
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/mutsa-delivery
EnvironmentFile=/opt/mutsa-delivery/.env
ExecStart=/usr/bin/java -Xms64m -Xmx192m -XX:+UseSerialGC -jar /opt/mutsa-delivery/app.jar
SuccessExitStatus=143
Restart=on-failure
RestartSec=10
MemoryHigh=260M
MemoryMax=350M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
